# config/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /home/server/src/config/uwsgi.ini --catch-exceptions

[program:nginx-app]
command = /usr/sbin/nginx
